# systemd/update-mirror.timer
[Unit]
Description=Update local mirrors

[Timer]
OnUnitActiveSec=20m

[Install]
WantedBy=timers.target
