# systemd/update-mirror-once-a-day.timer
[Unit]
Description=Update local mirrors (once a day)

[Timer]
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
