# systemd/update-mirror.timer
[Unit]
Description=Update local mirrors

[Timer]
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
